# airootfs/etc/skel/.config/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Install System
Exec=/usr/local/bin/start-calamares.sh
Icon=calamares
Terminal=false
