# airootfs/etc/systemd/system/calamares.service
[Unit]
Description=Calamares Installer

[Service]
ExecStart=/usr/bin/calamares
Restart=no

[Install]
WantedBy=graphical.target
